# mukund-dev/ansible/roles/jenkins/templates/jenkins.service
[Unit]
Description=Jenkins_Daemon

[Service]
ExecStart=/usr/bin/java -jar /usr/local/etc/jenkins.war


[Install]
WantedBy=multi-user.target
